# sigma-rules-intelligent/m365/entra_id/ALPHV-BlackCat/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: 10bb8029-69c9-484c-9771-4565b5230284
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: ALPHV-BlackCat
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.alphv-blackcat
