# sigma-rules-intelligent/m365/entra_id/Stealth Falcon/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 68546cc4-394a-4aa7-bfb6-c5ef8d0646f2
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Stealth Falcon
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.stealth_falcon
