# sigma-rules/google/workspace/inc_ransom/t1053.yml
title: Scheduled Task/Job - INC Ransom (workspace)
id: 6ed9d044-75ed-437e-98f3-70b869381208
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053 (Scheduled Task/Job).
  Associated threat actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053/
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_cron_jobs_creation_and_runtime.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053
  - actor.inc_ransom
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    process.args|contains: '-l'
    process.name|contains: 'running job'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
