# sigma-rules-intelligent/windows/sysmon/Stealth Falcon/T1053.005.yml
title: Scheduled Task (WINDOWS - sysmon)
id: af09f057-392e-4175-8eb2-587c1358491b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Stealth Falcon
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.stealth_falcon
