# sigma-rules-intelligent/windows/sysmon/BRONZE BUTLER/T1053.005.yml
title: Scheduled Task (WINDOWS - sysmon)
id: a34d3977-9d85-4e5c-8024-1d32909d96a3
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BRONZE BUTLER
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.bronze_butler
