# sigma-rules/m365/entra_id/cactus_ransomware/t1053.005.yml
title: Scheduled Task - Cactus Ransomware (entra_id)
id: d06ac857-617f-4021-87c0-10d62d7316f8
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Cactus Ransomware
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.cactus_ransomware
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
