# sigma-rules/macos/unified_logging/cactus_ransomware/t1053.005.yml
title: Scheduled Task - Cactus Ransomware (unified_logging)
id: 157bd35a-c7b6-489e-885b-ecd21b5d5acc
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Cactus Ransomware
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.cactus_ransomware
logsource:
  product: macos
  service: unified_logging
  category: process_creation
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
