# sigma-rules/windows/general/t1053.005.yml
title: Scheduled Task Executed From A Suspicious Location
id: 5dea569f-5dea-4dea-adea-5dea569f0000
status: stable
description: Detection rule for MITRE ATT&CK technique T1053.005
references:
  - https://github.com/SigmaHQ/sigma/blob/master/rules/windows/builtin/taskscheduler/win_taskscheduler_execution_from_susp_locations.yml
author: sigma (via IncidentBuddy)
date: 2026-01-28
tags:
  - attack.t1053.005
logsource:
  product: windows
  service: general
detection:
  selection:
    # Original sigma detection logic
    #     selection:
    #         EventID: 129 # Created Task Process
    #         Path|contains:
    #             - 'C:\Windows\Temp\'
    #             - '\AppData\Local\Temp\'
    #             - '\Desktop\'
    #             - '\Downloads\'
    #             - '\Users\Public\'
    #             - 'C:\Temp\'
    #     # If you experience FP. Uncomment the filter below and add the specific TaskName with the Program to it
    #     # filter:
    #     #     TaskName: '\Exact\Task\Name'
    #     #     Path: 'Exact\Path'
    #     condition: selection
    # 
  condition: selection
level: medium
