# sigma-rules-intelligent/linux/auditd/Brain Cipher/T1053.005.yml
title: Scheduled Task (LINUX - auditd)
id: 1189f919-a37b-4362-a8e1-783d182d3760
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.brain_cipher
