# sigma-rules-intelligent/macos/unified_logging/BQTLock/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: efdb135e-026b-4473-85f9-42c820a48d67
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BQTLock
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.bqtlock
