# sigma-rules-intelligent/windows/defender/Nitrogen/T1053.005.yml
title: Scheduled Task (WINDOWS - defender)
id: 7474533d-cf8b-473a-b34a-161c0bb3f693
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Nitrogen
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.nitrogen
