# sigma-rules-intelligent/windows/security/Nitrogen/T1053.005.yml
title: Scheduled Task (WINDOWS - security)
id: e2301af3-7057-4d53-bf70-e57af6eb0347
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Nitrogen
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    NewProcessName|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.nitrogen
