# sigma-rules-intelligent/google/workspace/ALPHV-BlackCat/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: d642f244-c29a-4c00-9638-efe97392b8d4
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: ALPHV-BlackCat
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.alphv-blackcat
