# sigma-rules-intelligent/macos/unified_logging/BianLian/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 8c84185c-72da-47e1-aae6-1eea9e69793b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BianLian
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    ProcessName|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.bianlian
