# sigma-rules-intelligent/google/workspace/Brain Cipher/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: b93eb818-80f6-4cbc-89e0-111e2ee54470
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.brain_cipher
