# sigma-rules-intelligent/google/workspace/Brain Cipher/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: ed1fe61b-077e-4ed7-9240-2b69026f8934
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.brain_cipher
