# sigma-rules-intelligent/m365/entra_id/APT41/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 98777e06-5e09-46ff-8392-f1daa4e44043
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: APT41
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    userPrincipalName: user@domain.com
    ipAddress|startswith:
      - '192.168'
      - '10.0'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.apt41
