# sigma-rules/google/workspace/nitrogen/t1053.yml
title: Scheduled Task/Job - Nitrogen (workspace)
id: 4d5c2d3e-37df-4a36-b1fb-569ee7fbf4a5
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053 (Scheduled Task/Job).
  Associated threat actor: Nitrogen
references:
  - https://attack.mitre.org/techniques/T1053/
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_cron_jobs_creation_and_runtime.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053
  - actor.nitrogen
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    process.args|contains: '-l'
    process.name|contains: 'running job'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
