# sigma-rules/google/workspace/nitrogen/t1053.005.yml
title: Scheduled Task - Nitrogen (workspace)
id: d8bf2883-ffc2-4cf2-996e-536bbb5d9fc3
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Nitrogen
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.nitrogen
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
