# sigma-rules-intelligent/windows/sysmon/BlackBasta/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: ec428ac5-b0aa-4ac1-b335-bfb6a6f5d3bb
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BlackBasta
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.blackbasta
