# sigma-rules-intelligent/m365/entra_id/COOSEA/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 9aa318b4-328b-4c85-9214-84b4ffe152eb
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: COOSEA
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.coosea
