# sigma-rules-intelligent/google/workspace/Cactus Ransomware/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: c67ea67f-e7f1-4358-a705-3c252bc6452b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Cactus Ransomware
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.cactus_ransomware
