# sigma-rules-intelligent/google/workspace/APT39/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: b9433fa8-7f3d-4ff8-a4e2-46ebe42a03a7
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: APT39
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.apt39
