# sigma-rules-intelligent/google/workspace/COOSEA/T1053.005.yml
title: Scheduled Task (GOOGLE)
id: 5ee5f476-db8b-4b6f-aa8f-14a72e1ab9bc
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
