# sigma-rules-intelligent/google/INC Ransom/T1053.yml
title: Scheduled Task/Job (GOOGLE)
id: b68af3fe-3c07-4661-8f5f-3ef6eded739d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
