# sigma-rules/google/workspace/qilin/t1053.005.yml
title: Scheduled Task - Qilin (workspace)
id: ae28aa85-5d9c-4a86-9e1e-b0471d04b301
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Qilin
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.qilin
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
