# sigma-rules-intelligent/google/workspace/BianLian/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: 1d677bd3-f78f-4920-97c3-8addbd935c56
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BianLian
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.bianlian
