# sigma-rules-intelligent/google/workspace/BianLian/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 3fa8309e-81f2-4085-9b41-3cbd43f0f672
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BianLian
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.bianlian
