# sigma-rules/m365/entra_id/lynx/t1053.005.yml
title: Scheduled Task - Lynx (entra_id)
id: ab47258e-626f-49f6-b4d7-40537afb6524
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Lynx
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.lynx
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
