# sigma-rules/google/workspace/cactus_ransomware/t1053.005.yml
title: Scheduled Task - Cactus Ransomware (workspace)
id: b4610968-349b-40f0-8bee-92a4a5b727f0
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Cactus Ransomware
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.cactus_ransomware
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
