# sigma-rules-intelligent/m365/entra_id/Brain Cipher/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: 5b848806-720c-4acb-92b7-259024d810f6
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.brain_cipher
