# sigma-rules-intelligent/windows/sysmon/GALLIUM/T1053.005.yml
title: Scheduled Task (WINDOWS - sysmon)
id: 1da59c3a-79de-469e-8b6c-02fdcfeffdd6
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: GALLIUM
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    UtcTime|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.gallium
