# sigma-rules-intelligent/windows/defender/Qilin/T1053.005.yml
title: Scheduled Task (WINDOWS - defender)
id: f17e0cb3-2073-4ba8-bb91-e4d0ef3c4092
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Qilin
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.qilin
