# sigma-rules-intelligent/windows/defender/Sarcoma/T1053.yml
title: Scheduled Task/Job (WINDOWS - defender)
id: 65e44cc6-acfa-40f5-8ff7-1f8dc7608abe
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Sarcoma
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.sarcoma
