# sigma-rules-intelligent/windows/defender/Sarcoma/T1053.yml
title: Scheduled Task/Job (WINDOWS - defender)
id: 34237d3a-14b8-445d-a659-91c5b620ea1b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Sarcoma
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.sarcoma
