# sigma-rules-intelligent/google/BianLian/T1053.yml
title: Scheduled Task/Job (GOOGLE)
id: 298c724f-cec6-4202-8030-46b11d27b75a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
