# sigma-rules-intelligent/m365/entra_id/Rancor/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 54749441-9d9b-4828-9ae5-d9c171667953
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Rancor
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.rancor
