# sigma-rules-intelligent/m365/entra_id/Fox Kitten/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 7bf50d55-816a-4e04-8347-4ed0d865ddd1
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Fox Kitten
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    userPrincipalName: user@domain.com
    ipAddress|startswith:
      - '192.168'
      - '10.0'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.fox_kitten
