# sigma-rules-intelligent/windows/defender/Stormous/T1053.yml
title: Scheduled Task/Job (WINDOWS - defender)
id: 2aba4ab6-68ae-480b-a561-540d3518264a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Stormous
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    InitiatingProcessFileName|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    ProcessCommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.stormous
