# sigma-rules/m365/entra_id/alphv-blackcat/t1053.yml
title: Scheduled Task/Job - ALPHV-BlackCat (entra_id)
id: 10556a63-46e6-4eff-8323-c3449443686e
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053 (Scheduled Task/Job).
  Associated threat actor: ALPHV-BlackCat
references:
  - https://attack.mitre.org/techniques/T1053/
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_cron_jobs_creation_and_runtime.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053
  - actor.alphv-blackcat
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    process.args|contains: '-l'
    process.name|contains: 'running job'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
