# sigma-rules/m365/entra_id/kazu/t1543.002.yml
title: Systemd Service - Kazu (entra_id)
id: 9f0a0ff9-d3c4-469c-b600-5fe59bef101c
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1543.002 (Systemd Service).
  Associated threat actor: Kazu
references:
  - https://attack.mitre.org/techniques/T1543/002/
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_etc_file_creation.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1543.002
  - actor.kazu
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    file.path|contains:
      - '/etc/ld.so.conf.d/*'
      - '/usr/lib/systemd/system/*'
    process.name|contains:
      - 'chef-client'
      - 'rhn_check'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
