# sigma-rules-intelligent/windows/sysmon/Stormous/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: bb62fa7e-b363-4fa2-8354-ff9af290a13c
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Stormous
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.stormous
