# sigma-rules/linux/auditd/cephalus/t1053.005.yml
title: Scheduled Task - Cephalus (auditd)
id: 6d4bd73c-7d00-4df1-b33a-0d8a9fe58fe8
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Cephalus
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.cephalus
logsource:
  product: linux
  service: auditd
  category: process_creation
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
