# sigma-rules-intelligent/m365/entra_id/MuddyWater/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 2e4fc4e2-fcf6-414f-9c13-bbda7b837491
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: MuddyWater
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.muddywater
