# sigma-rules-intelligent/google/workspace/SafePay/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: fbc787dc-e406-400b-b276-34f0841ed26e
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: SafePay
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.safepay
