# sigma-rules-intelligent/macos/unified_logging/Wizard Spider/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 41cdbecf-2ac3-4e66-8b5f-c2327a13eb98
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Wizard Spider
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    exe|contains:
      - '/bin/bash'
      - '/bin/sh'
      - '/usr/bin/python'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.wizard_spider
