# sigma-rules-intelligent/macos/unified_logging/Stormous/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: 0886b487-1733-48fc-8048-f6299cf4410b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Stormous
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.stormous
