# sigma-rules-intelligent/windows/sysmon/GLOBAL/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: 599120ac-7a24-4cee-997b-f2858e82daa0
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: GLOBAL
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.global
