# sigma-rules/macos/unified_logging/t1053.003.yml
title: Scheduled Cron Task/Job - MacOs
id: 4bc4d36b-06ba-402a-a27a-8a31b6b795fe
status: experimental
description: |
  Detection for MITRE ATT&CK technique T1053.003.
  Original source: sigma
references:
  - https://attack.mitre.org/techniques/T1053/003/
  - https://github.com/SigmaHQ/sigma/blob/master/rules/macos/process_creation/proc_creation_macos_schedule_task_job_cron.yml
author: GapMATRIX (sourced from sigma)
date: 2026-01-28
tags:
  - attack.t1053.003
logsource:
  product: macos
  service: unified_logging
  category: process_creation
detection:
      selection:
          Image|endswith: '/crontab'
          CommandLine|contains: '/tmp/'
      condition: selection
  
  condition: selection
falsepositives:
  - Legitimate administrative activity
level: medium
