# sigma-rules/macos/unified_logging/t1053.yml
title: Suspicious CronTab Creation or Modification
id: 9439a4a2-9c3c-4e79-be16-86a54b83f8bf
status: experimental
description: |
  Detection for MITRE ATT&CK technique T1053.
  Original source: elastic
references:
  - https://attack.mitre.org/techniques/T1053/
  - https://github.com/elastic/detection-rules/blob/main/rules/macos/persistence_crontab_creation.toml
author: GapMATRIX (sourced from elastic)
date: 2026-01-28
tags:
  - attack.t1053
logsource:
  product: macos
  service: unified_logging
  category: process_creation
detection:
  selection:
    placeholder: true
  condition: selection
falsepositives:
  - Legitimate administrative activity
level: medium
