# sigma-rules/macos/general/t1053.yml
title: Suspicious CronTab Creation or Modification
id: 658e3772-658e-458e-a58e-658e37720000
status: stable
description: Detection rule for MITRE ATT&CK technique T1053
references:
  - https://github.com/elastic/detection-rules/blob/main/rules/macos/persistence_crontab_creation.toml
author: elastic (via IncidentBuddy)
date: 2026-02-08
tags:
  - attack.t1053
logsource:
  product: macos
  service: general
detection:
  selection:
    # Original elastic detection logic
    # file where host.os.type == "macos" and event.type != "deletion" and process.name != null and
    #   file.path like "/private/var/at/tabs/*" and not process.executable == "/usr/bin/crontab"
  condition: selection
level: medium
