# sigma-rules/macos/general/t1543.yml
title: Deprecated - LaunchDaemon Creation or Modification and Immediate Loading
id: 04355339-0435-4435-a435-043553390000
status: stable
description: Detection rule for MITRE ATT&CK technique T1543
references:
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_creation_modif_launch_deamon_sequence.toml
author: elastic (via IncidentBuddy)
date: 2026-02-08
tags:
  - attack.t1543
logsource:
  product: macos
  service: general
detection:
  selection:
    # Original elastic detection logic
    # sequence by host.id with maxspan=1m
    #  [file where host.os.type == "macos" and event.type != "deletion" and file.path : ("/System/Library/LaunchDaemons/*", "/Library/LaunchDaemons/*")]
    #  [process where host.os.type == "macos" and event.type in ("start", "process_started") and process.name == "launchctl" and process.args == "load"]
  condition: selection
level: medium
