# sigma-rules-intelligent/google/workspace/BQTLock/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: d95ad456-a65f-4c2a-abc0-17f6e5f8c436
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BQTLock
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.bqtlock
