# sigma-rules/macos/unified_logging/lynx/t1053.005.yml
title: Scheduled Task - Lynx (unified_logging)
id: 000c2891-0001-43a4-a6d8-83de5e3b514a
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Lynx
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.lynx
logsource:
  product: macos
  service: unified_logging
  category: process_creation
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
