# sigma-rules-intelligent/m365/entra_id/APT-C-36/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: d466aa4d-1a7b-4b0f-a97e-1a1aaf12f89d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: APT-C-36
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    userPrincipalName: user@domain.com
    ipAddress|startswith:
      - '192.168'
      - '10.0'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.apt-c-36
