# sigma-rules-intelligent/windows/defender/Confucius/T1053.005.yml
title: Scheduled Task (WINDOWS - defender)
id: 3970fb49-9bd1-43ff-a613-8daf0c3d7e43
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Confucius
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.confucius
