# sigma-rules-intelligent/google/workspace/Sarcoma/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: ea4b8cd0-f2e1-4820-8e9c-4d30efda27f0
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Sarcoma
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.sarcoma
