# sigma-rules-intelligent/windows/security/BlackByte/T1053.yml
title: Scheduled Task/Job (WINDOWS - security)
id: 2cd93245-b97f-4d8b-a47d-72cdc32b2bc6
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BlackByte
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.blackbyte
