# sigma-rules-intelligent/windows/security/BlackByte/T1053.005.yml
title: Scheduled Task (WINDOWS - security)
id: 0305d1a2-5e8d-40db-8b48-ef3c1afd2ebc
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BlackByte
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    NewProcessName|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.blackbyte
