# sigma-rules-intelligent/m365/entra_id/BianLian/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: a1b45d23-ad3e-41e7-b142-ca76dc9e1d4e
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BianLian
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.bianlian
