# sigma-rules-intelligent/google/workspace/Rancor/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: f7ea897e-084e-471f-9fd1-e7154e26df57
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Rancor
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.rancor
