# sigma-rules-intelligent/m365/entra_id/Lazarus Group/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 681c1db4-62c0-44f3-9e31-7be56e865002
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Lazarus Group
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.lazarus_group
