# sigma-rules-intelligent/google/workspace/INC Ransom/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 10b77932-d937-443f-a88f-96a4e5927d5e
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.inc_ransom
