# sigma-rules-intelligent/google/workspace/INC Ransom/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: d00f5827-68db-4bd4-9310-1589352c4693
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.inc_ransom
