# sigma-rules-intelligent/google/workspace/INC Ransom/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: b66c4426-c6d3-4f23-a0a8-1288540f6a94
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.inc_ransom
