# sigma-rules-intelligent/google/workspace/Confucius/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 5b4a7a71-b746-423f-9075-11dda80a3a93
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Confucius
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.confucius
