# sigma-rules-intelligent/windows/Yashma-Ator-Vietnamita/T1053.yml
title: Scheduled Task/Job (WINDOWS)
id: 7d19fca4-cc81-45d6-9735-f5946b9b23a8
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
