# sigma-rules-intelligent/macos/unified_logging/APT39/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: d01a4fc3-5231-4814-b7c5-83b9da294765
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: APT39
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    ProcessName|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.apt39
