# sigma-rules-intelligent/macos/unified_logging/Cactus Ransomware/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: b217629f-e0e4-4ef6-9dbf-36c6b12fd9a8
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Cactus Ransomware
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.cactus_ransomware
