# sigma-rules-intelligent/macos/unified_logging/Cactus Ransomware/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 788c30e1-68c3-4acd-878a-00ca36f3535b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Cactus Ransomware
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd'
      - 'powershell'
    CommandLine|contains:
      - 'encoded'
      - 'bypass'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.cactus_ransomware
