# sigma-rules-intelligent/linux/auditd/Spectra/T1053.yml
title: Scheduled Task/Job (LINUX - auditd)
id: e5382acb-5c2f-4b70-ad79-008b1b28b16d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Spectra
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    exe|contains:
      - '/bin/bash'
      - '/bin/sh'
      - '/usr/bin/python'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.spectra
