# sigma-rules-intelligent/macos/unified_logging/OilRig/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 7515cf6c-b257-41e5-ac16-761ca24986ee
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: OilRig
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd'
      - 'powershell'
    CommandLine|contains:
      - 'encoded'
      - 'bypass'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.oilrig
