# sigma-rules-intelligent/macos/unified_logging/Rancor/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: fa3b92fb-e0d6-40be-a82a-aa72fc22709a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Rancor
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    ProcessName|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.rancor
