# sigma-rules-intelligent/google/workspace/GLOBAL/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: 2c14d7d8-2cf2-492b-858d-6411a0de2b3d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: GLOBAL
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.global
