# sigma-rules-intelligent/google/workspace/GLOBAL/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: c260aba2-7dff-463f-98f4-37ba437cc716
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: GLOBAL
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.global
