# sigma-rules-intelligent/m365/entra_id/Spectra/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: 87d1abfc-e1dd-4a74-80a2-ea8c3d29275d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Spectra
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.spectra
