# sigma-rules-intelligent/m365/entra_id/Spectra/T1053.yml
title: Scheduled Task/Job (M365)
id: f2eb074a-e07a-46a3-ab5b-9c2b4ddfcf4b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
