# sigma-rules-intelligent/m365/entra_id/Vect RaaS/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: d11db3f6-d506-4fd4-b3f5-167e1f8456e6
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Vect RaaS
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    userPrincipalName: user@domain.com
    ipAddress|startswith:
      - '192.168'
      - '10.0'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.vect_raas
