# sigma-rules-intelligent/google/workspace/FIN10/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: c77263be-74c2-4a32-9494-310fbb1e82fc
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: FIN10
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.fin10
