# sigma-rules-intelligent/m365/entra_id/Qilin/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: 88ba8f9b-1894-41ed-a9f0-700925d215a7
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Qilin
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.qilin
