# sigma-rules-intelligent/google/workspace/Fox Kitten/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: ecc9210a-5bc0-4b02-a94e-56f313b11926
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Fox Kitten
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.fox_kitten
