# sigma-rules-intelligent/macos/unified_logging/Vect RaaS/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: f96af375-b264-4322-aa07-5df670113c3e
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Vect RaaS
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd'
      - 'powershell'
    CommandLine|contains:
      - 'encoded'
      - 'bypass'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.vect_raas
