# sigma-rules/windows/defender/inc_ransom/t1053.005.yml
title: Scheduled Task - INC Ransom (defender)
id: dd7c6810-dc9a-4d4e-8034-08938508bc36
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://github.com/elastic/detection-rules/blob/main/rules/windows/execution_scheduled_task_powershell_source.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.inc_ransom
logsource:
  product: windows
  service: defender
  category: process_creation
detection:
  selection:
    event.action|contains: 'Image loaded*'
    dll.name|contains: 'taskschd.dll'
    file.name|contains: 'taskschd.dll'
    process.name|contains:
      - 'powershell.exe'
      - 'powershell_ise.exe'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
