# sigma-rules-intelligent/windows/security/Stormous/T1053.yml
title: Scheduled Task/Job (WINDOWS - security)
id: b357da6c-f24a-4fb4-8c95-030972855ab2
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Stormous
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.stormous
