# sigma-rules-intelligent/m365/entra_id/Lynx/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: 19e9dff3-4ddc-4ff8-8864-f0c6127ee2de
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Lynx
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.lynx
