# sigma-rules/m365/entra_id/bianlian/t1053.005.yml
title: Scheduled Task - BianLian (entra_id)
id: 0dc52b5f-0ed4-4bf9-a665-ff5027b5d823
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: BianLian
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.bianlian
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
