# sigma-rules-intelligent/macos/unified_logging/ALPHV-BlackCat/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: 7e4d1ced-8c5c-4e01-a9f2-1a59a42f79cb
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: ALPHV-BlackCat
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.alphv-blackcat
