# sigma-rules-intelligent/macos/unified_logging/ALPHV-BlackCat/T1053.003.yml
title: Cron (MACOS - unified_logging)
id: 5a5d48b6-6ded-4f81-901b-5741f1aecc1f
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.003.
  Threat Actor: ALPHV-BlackCat
references:
  - https://attack.mitre.org/techniques/T1053.003/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    Image|endswith|contains:
      - '/crontab'
    CommandLine|contains|contains:
      - '/tmp/'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.003
  - actor.alphv-blackcat
