# sigma-rules-intelligent/google/workspace/Wizard Spider/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 849bc62d-1669-49e4-80b2-9db2c2b0a752
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Wizard Spider
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.wizard_spider
