# sigma-rules-intelligent/google/workspace/BlackByte/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 023110a3-a0b5-4035-b456-4a4429b80ff7
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: BlackByte
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.blackbyte
