# sigma-rules/google/workspace/brain_cipher/t1053.005.yml
title: Scheduled Task - Brain Cipher (workspace)
id: 4404435e-a3ee-45a4-9c47-d4232aaa37fd
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.brain_cipher
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
