# sigma-rules-intelligent/windows/defender/GLOBAL/T1053.yml
title: Scheduled Task/Job (WINDOWS)
id: 523e5e99-33fa-4e98-b6ec-6f0c701b9a6f
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
