# sigma-rules-intelligent/m365/entra_id/Yashma-Ator-Vietnamita/T1053.yml
title: Scheduled Task/Job (M365 - entra_id)
id: 01b49b65-35a2-4d32-aea2-53b40d9b591b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Yashma-Ator-Vietnamita
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.yashma-ator-vietnamita
