# sigma-rules-intelligent/google/workspace/Sandworm Team/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 38a5f9e0-35ad-40cf-a573-f33a41e68fb5
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Sandworm Team
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.sandworm_team
