# sigma-rules-intelligent/google/workspace/Magic Hound/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 8540d2a5-d955-4686-9068-89a2e72a7028
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Magic Hound
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.magic_hound
