# sigma-rules-intelligent/windows/security/SafePay/T1053.005.yml
title: Scheduled Task (WINDOWS)
id: 16cc9638-3827-40ac-bffc-078213603a25
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
