# sigma-rules-intelligent/macos/unified_logging/FIN10/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 246234f9-4be9-46d0-9ae5-15d608be1cc6
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: FIN10
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd'
      - 'powershell'
    CommandLine|contains:
      - 'encoded'
      - 'bypass'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.fin10
