# sigma-rules-intelligent/windows/sysmon/Dragonfly/T1053.005.yml
title: Scheduled Task (WINDOWS - sysmon)
id: 3efd4ac3-7879-4c99-8d3a-7281d4af1014
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Dragonfly
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.dragonfly
