# sigma-rules-intelligent/windows/sysmon/BlackByte/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: d30e3e1f-5a83-41fd-9e66-6c50446e85a3
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BlackByte
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    Image|endswith|contains:
      - '\crackmapexec.exe'
    CommandLine|contains|contains:
      - ' -t schtask -c '
      - ' -t startupfolder -c '
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.blackbyte
