# sigma-rules-intelligent/windows/sysmon/BlackByte/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: d63d01a8-6961-46f5-a0b6-e08785627e51
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: BlackByte
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    UtcTime|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.blackbyte
