# sigma-rules-intelligent/m365/entra_id/Sandworm Team/T1053.005.yml
title: Scheduled Task (M365 - entra_id)
id: d31bd38a-0d73-4752-9dd6-582d715fa4ec
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Sandworm Team
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: m365
  service: entra_id
  category: authentication

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.sandworm_team
