# sigma-rules-intelligent/macos/unified_logging/Brain Cipher/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: f302c7b2-c155-42e9-9f0a-90b6a4fac17d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Brain Cipher
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.brain_cipher
