# sigma-rules-intelligent/macos/unified_logging/Generic/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: 36a51c87-8d69-4b2f-84e1-2a6f1fba3d10
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Generic
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.generic
