# sigma-rules/m365/entra_id/cephalus/t1053.005.yml
title: Scheduled Task - Cephalus (entra_id)
id: e53903c9-a3cc-4b8b-b6cc-a4b6da03943b
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Cephalus
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.cephalus
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
