# sigma-rules-intelligent/windows/defender/COOSEA/T1053.yml
title: Scheduled Task/Job (WINDOWS - defender)
id: 41a322e2-fa35-44c0-bb3a-0250041d7077
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: COOSEA
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    Image|endswith|contains:
      - '\crackmapexec.exe'
    CommandLine|contains|contains:
      - ' -t schtask -c '
      - ' -t startupfolder -c '
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.coosea
