# sigma-rules-intelligent/windows/sysmon/Generic/T1053.yml
title: Scheduled Task/Job (WINDOWS - sysmon)
id: 7f3b9d7e-5e16-47a4-9912-53add94cfdea
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Generic
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: sysmon
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'cmd.exe'
      - 'powershell.exe'
      - 'bash'
      - 'sh'
    CommandLine|contains:
      - '-enc'
      - 'bypass'
      - 'encoded'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.generic
