# sigma-rules-intelligent/google/workspace/Machete/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 5c63aa7e-1ba3-412d-a940-40feb8929b9b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Machete
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.machete
