# sigma-rules/m365/entra_id/brain_cipher/t1543.yml
title: Create or Modify System Process - Brain Cipher (entra_id)
id: 32163cfb-7dac-4358-8e0c-cc804bc3e765
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1543 (Create or Modify System Process).
  Associated threat actor: Brain Cipher
references:
  - https://attack.mitre.org/techniques/T1543/
  - https://github.com/elastic/detection-rules/blob/main/rules/_deprecated/persistence_etc_file_creation.toml
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1543
  - actor.brain_cipher
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    file.path|contains:
      - '/etc/ld.so.conf.d/*'
      - '/usr/lib/systemd/system/*'
    process.name|contains:
      - 'chef-client'
      - 'rhn_check'
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
