# sigma-rules-intelligent/linux/Qilin/T1053.005.yml
title: Scheduled Task (LINUX)
id: 1d8351a4-f965-4f22-8434-9e699871b679
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    Image|contains: process
    CommandLine|contains: cmd
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
