# sigma-rules/linux/general/t1053.002.yml
title: Scheduled Task/Job At
id: 6a84d433-6a84-4a84-aa84-6a84d4330000
status: stable
description: Detection rule for MITRE ATT&CK technique T1053.002
references:
  - https://github.com/SigmaHQ/sigma/blob/master/rules/linux/process_creation/proc_creation_lnx_at_command.yml
author: sigma (via IncidentBuddy)
date: 2026-02-01
tags:
  - attack.t1053.002
logsource:
  product: linux
  service: general
detection:
  selection:
    # Original sigma detection logic
    #     selection:
    #         Image|endswith:
    #             - '/at'
    #             - '/atd'
    #     condition: selection
    # 
  condition: selection
level: medium
