# sigma-rules/linux/general/t1053.003.yml
title: Modifying Crontab
id: 127ed1b4-127e-427e-a27e-127ed1b40000
status: stable
description: Detection rule for MITRE ATT&CK technique T1053.003
references:
  - https://github.com/SigmaHQ/sigma/blob/master/rules/linux/builtin/cron/lnx_cron_crontab_file_modification.yml
author: sigma (via IncidentBuddy)
date: 2026-02-01
tags:
  - attack.t1053.003
logsource:
  product: linux
  service: general
detection:
  selection:
    # Original sigma detection logic
    #     keywords:
    #         - 'REPLACE'
    #     condition: keywords
    # 
  condition: selection
level: medium
