# sigma-rules-intelligent/windows/security/8base/T1053.yml
title: Scheduled Task/Job (WINDOWS - security)
id: 1ed8496e-4c49-4fd4-aacb-8c4ac536c1dc
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: 8base
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.8base
