# sigma-rules-intelligent/linux/auditd/Magic Hound/T1053.005.yml
title: Scheduled Task (LINUX - auditd)
id: e639f47b-e13b-432c-a6e6-9adf891362c0
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Magic Hound
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    exe|contains:
      - '/bin/bash'
      - '/bin/sh'
      - '/usr/bin/python'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.magic_hound
