# sigma-rules-intelligent/windows/security/INC Ransom/T1053.005.yml
title: Scheduled Task (WINDOWS - security)
id: d0cc7400-7349-41ed-942b-a9a4200557f5
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: security
  category: process_creation

detection:
  selection:
    Image|contains:
      - 'suspicious'
      - 'malicious'
      - 'attack'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.inc_ransom
