# sigma-rules-intelligent/macos/unified_logging/FIN8/T1053.005.yml
title: Scheduled Task (MACOS - unified_logging)
id: 74f3fcde-9b7f-4b7b-a859-86dcdf474e2a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: FIN8
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    exe|contains:
      - '/bin/bash'
      - '/bin/sh'
      - '/usr/bin/python'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.fin8
