# sigma-rules/google/workspace/lynx/t1053.005.yml
title: Scheduled Task - Lynx (workspace)
id: de3e204a-9a61-4b24-b7ee-3ef4e871b3f1
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: Lynx
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.lynx
logsource:
  product: google
  service: workspace
  category: admin_activity
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
