# sigma-rules-intelligent/google/workspace/Vect RaaS/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: 91761c52-b172-4110-8ed8-0568ef6e9f8a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Vect RaaS
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.vect_raas
