# sigma-rules-intelligent/linux/auditd/INC Ransom/T1053.005.yml
title: Scheduled Task (LINUX - auditd)
id: bdb6ebb4-04f2-468c-a1ae-9af98ab0e6d0
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: INC Ransom
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.inc_ransom
