# sigma-rules-intelligent/linux/auditd/INC Ransom/T1053.yml
title: Scheduled Task/Job (LINUX - auditd)
id: c132ac3b-af50-4977-8ba1-ebfc900b3cdd
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: INC Ransom
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: linux
  service: auditd
  category: process_creation

detection:
  selection:
    exe|contains:
      - '/bin/bash'
      - '/bin/sh'
      - '/usr/bin/python'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.inc_ransom
