# sigma-rules-intelligent/macos/unified_logging/INC Ransom/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: 8404a9da-481a-4358-92a7-6b57e7391d11
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: INC Ransom
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.inc_ransom
