# sigma-rules-intelligent/google/workspace/Generic/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: 5d39444e-fbcf-422b-9d13-2ad99e03a83a
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Generic
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    event|contains: action
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.generic
