# sigma-rules/windows/powershell/t1053.005.yml
title: Suspicious Scheduled Task Creation
id: 215f0269-215f-415f-a15f-215f02690000
status: stable
description: Detection rule for MITRE ATT&CK technique T1053.005
references:
  - https://github.com/SigmaHQ/sigma/blob/master/rules/windows/builtin/security/win_security_susp_scheduled_task_creation.yml
author: sigma (via IncidentBuddy)
date: 2026-02-08
tags:
  - attack.t1053.005
logsource:
  product: windows
  service: powershell
detection:
  selection:
    # Original sigma detection logic
    #     selection_eid:
    #         EventID: 4698
    #     selection_paths:
    #         TaskContent|contains:
    #             - '\AppData\Local\Temp\'
    #             - '\AppData\Roaming\'
    #             - '\Users\Public\'
    #             - '\WINDOWS\Temp\'
    #             - 'C:\Temp\'
    #             - '\Desktop\'
    #             - '\Downloads\'
    #             - '\Temporary Internet'
    #             - 'C:\ProgramData\'
    #             - 'C:\Perflogs\'
    #     selection_commands:
    #         TaskContent|contains:
    #             - 'regsvr32'
    #             - 'rundll32'
    #             - 'cmd.exe</Command>'
    #             - 'cmd</Command>'
    #             - '<Arguments>/c '
    #             - '<Arguments>/k '
    #             - '<Arguments>/r '
    #             - 'powershell'
    #             - 'pwsh'
    #             - 'mshta'
    #             - 'wscript'
    #             - 'cscript'
    #             - 'certutil'
    #             - 'bitsadmin'
    #             - 'bash.exe'
    #             - 'bash '
    #             - 'scrcons'
    #             - 'wmic '
    #             - 'wmic.exe'
    #             - 'forfiles'
    #             - 'scriptrunner'
    #             - 'hh.exe'
    #     condition: all of selection_*
    # 
  condition: selection
level: medium
