# sigma-rules-intelligent/google/workspace/8base/T1053.yml
title: Scheduled Task/Job (GOOGLE - workspace)
id: 794f90a6-222f-4f30-ad99-424b12e0c8f1
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: 8base
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.8base
