# sigma-rules-intelligent/macos/unified_logging/COOSEA/T1053.yml
title: Scheduled Task/Job (MACOS - unified_logging)
id: 5d377f6b-343e-4efd-88aa-900614500a4b
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: COOSEA
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: macos
  service: unified_logging
  category: process_creation

detection:
  selection:
    placeholder: true
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.coosea
