# sigma-rules-intelligent/windows/defender/Yashma-Ator-Vietnamita/T1053.yml
title: Scheduled Task/Job (WINDOWS - defender)
id: 20f11364-daf2-45af-9d83-f5b68fa4911d
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.
  Threat Actor: Yashma-Ator-Vietnamita
  STATUS: Placeholder - waiting for detection indicators to be populated
references:
  - https://attack.mitre.org/techniques/T1053/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: windows
  service: defender
  category: process_creation

detection:
  selection:
    placeholder|contains:
      - 'PLACEHOLDER - Analyst to populate indicators'
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053
  - actor.yashma-ator-vietnamita
