# sigma-rules-intelligent/google/workspace/Mustang Panda/T1053.005.yml
title: Scheduled Task (GOOGLE - workspace)
id: 5ecf48f7-4618-4e7b-89d6-5b42671a9cce
description: >
  Detects techniques consistent with MITRE ATT&CK technique T1053.005.
  Threat Actor: Mustang Panda
references:
  - https://attack.mitre.org/techniques/T1053.005/
author: GapMATRIX Intelligent Sigma Generator
date: 2026-01-28
status: experimental
severity: medium

logsource:
  product: google
  service: workspace
  category: admin_activity

detection:
  selection:
    actor.email: admin@example.com
    events.type: admin_role_assigned
  condition: selection

falsepositives:
  - Legitimate system administration activity
  - Authorized security testing

tags:
  - attack.T1053.005
  - actor.mustang_panda
