# sigma-rules/linux/auditd/coosea/t1053.005.yml
title: Scheduled Task - COOSEA (auditd)
id: 1046ce3e-2d73-40fd-b199-c22384a8d614
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: COOSEA
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.coosea
logsource:
  product: linux
  service: auditd
  category: process_creation
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
