# sigma-rules/m365/entra_id/coosea/t1053.005.yml
title: Scheduled Task - COOSEA (entra_id)
id: c99ad739-5e1e-4951-975a-3523f3a5ac21
status: experimental
description: |
  Detects activity consistent with MITRE ATT&CK technique T1053.005 (Scheduled Task).
  Associated threat actor: COOSEA
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: GapMATRIX Sigma Generator
date: 2026-01-28
tags:
  - attack.t1053.005
  - actor.coosea
logsource:
  product: m365
  service: entra_id
  category: authentication
detection:
  selection:
    placeholder: true  # No detection logic available
  condition: selection
falsepositives:
  - Legitimate administrative activity
  - Authorised security testing
level: medium
